# crontab
*/5 * * * * /bin/bash /root/temper/temper.cron.sh > /root/temper/temper.cron.log 2> /root/temper/temper.cron.err
